# etc/crontabs/root
*/5 * * * * /usr/bin/thermalProtection
*/2 * * * * /sbin/trafficmon.lua
0 * * * * kill -14 `pgrep -f 'odhcp6c.*eth5'`
